# fetcher/etc/crontab
# run feed fetcher every five minutes
*/5 * * * * /opt/reader/bin/run --name fetcher