# fetcher/etc/crontab
# run feed fetcher every five minutes
*/5 * * * * reader /opt/reader/bin/run --name fetcher
