# parser/etc/crontab
# run feed parser every five minutes
*/5 * * * * reader sleep 3m && /opt/reader/bin/run --name parser
